# configs/vllm/vllm-switcher.service
[Unit]
Description=vLLM Model Switcher API
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/bin/vllm-switcher-api.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
