# configs/vllm/vllm-dan-qwen.service
[Unit]
Description=vLLM OpenAI-compatible API - Dan-Qwen3.5-4B
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/vllm-env
Environment=CUDA_VISIBLE_DEVICES=0
Environment=HF_HUB_OFFLINE=1
Environment=VLLM_WORKER_MULTIPROC_METHOD=spawn
ExecStart=/opt/vllm-env/bin/vllm serve /opt/models/dan-qwen3.5-4b     --host 0.0.0.0     --port 8000     --served-model-name dan-qwen3.5-4b     --dtype float16     --gpu-memory-utilization 0.95     --max-model-len 4096     --max-num-seqs 8     --trust-remote-code     --enforce-eager
Restart=on-failure
RestartSec=30
TimeoutStartSec=300
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vllm-dan-qwen

[Install]
WantedBy=multi-user.target
