# configs/vllm/vllm-badmistral.service
[Unit]
Description=vLLM OpenAI-compatible API - BADMISTRAL-1.5B
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/vllm-env
Environment="CUDA_VISIBLE_DEVICES=0"
Environment="HF_HUB_OFFLINE=1"
Environment="VLLM_WORKER_MULTIPROC_METHOD=spawn"
ExecStart=/opt/vllm-env/bin/vllm serve /opt/models/badmistral-1.5b \
    --host 0.0.0.0 \
    --port 8000 \
    --served-model-name badmistral-1.5b \
    --dtype float16 \
    --gpu-memory-utilization 0.90 \
    --max-model-len 4096 \
    --max-num-seqs 16 \
    --trust-remote-code
Restart=on-failure
RestartSec=30
TimeoutStartSec=300
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vllm-badmistral

[Install]
WantedBy=multi-user.target
